# configs/supervisor-app.conf
; Supervisor Configuration Example
; انسخ هذا الملف إلى /etc/supervisor/conf.d/ وعدّل حسب تطبيقك
;
; sudo cp supervisor-app.conf /etc/supervisor/conf.d/myapp.conf
; sudo supervisorctl reread
; sudo supervisorctl update

[program:myapp]
; الأمر لتشغيل التطبيق
command=/usr/bin/node /home/ubuntu/myapp/server.js

; مجلد العمل
directory=/home/ubuntu/myapp

; المستخدم
user=ubuntu

; بدء تلقائي
autostart=true

; إعادة تشغيل تلقائي عند الفشل
autorestart=true

; عدد محاولات إعادة التشغيل
startretries=10

; الانتظار قبل اعتبار البدء ناجحاً (ثواني)
startsecs=5

; إيقاف بـ SIGTERM
stopsignal=TERM

; الانتظار قبل SIGKILL (ثواني)
stopwaitsecs=10

; ملفات السجل
stdout_logfile=/var/log/myapp/stdout.log
stderr_logfile=/var/log/myapp/stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile_backups=5

; متغيرات البيئة
environment=
    NODE_ENV="production",
    PORT="3000"

; أولوية البدء (أقل = يبدأ أولاً)
priority=100
